# trunk/spells/laptop-mode-tools
   VERSION=( "1.05" "1.04" )
  CATEGORY="disk utility"
 ATTRIBUTE="console"
    SOURCE="downloads/laptop-mode-tools_$VERSION"
       URL="http://www.xs4all.nl/~bsamwel/laptop_mode/tools"
  HOMEPAGE="http://www.xs4all.nl/~bsamwel/laptop_mode/tools"
       REQ="acpid"
  ESTIMATE="5"
      DESC="$SPELL provides tools for the kernel laptop mode.

Laptop mode (also known as laptopmode, laptop-mode and 
laptop_mode, for search engine purposes :) ) is a kernel 
mode that allows you to extend the battery life of your 
laptop. It does this by intelligently grouping write activity 
on your disks, so that only reads of uncached data result in a 
disk spinup. It causes a significant improvement in battery 
life (for usage patterns that allow it)."

build(){
sed  -i  "s:/usr/local:/usr:" install.sh

sed  -i  "2i### BEGIN INIT INFO
2i# Provides: laptop-mode
2i# Required-Start: \$local_fs
2i# Required-Stop:  \$local_fs
2i# Default-Start: S
2i# Default-Stop: 
2i# Short-Description: laptop-mode allows to extend the battery life 
2i# Description: laptop-mode allows to extend the battery life by
2i#     grouping disk write activity, so that only reads of uncached
2i#     data result in a disk spinup.
2i### END INIT INFO"  \
etc/init.d/laptop-mode

 prepare_install	&&
 ./install.sh
}

post_install() {
 /lib/lsb/install_initd  /etc/init.d/laptop-mode
 true
}
