# trunk/spells/qpopper.d/xinetd.d/qpopper
service pop3
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/sbin/popper
	server_args	= -f /etc/qpopper.conf
	per_source	= 4
	log_type	= FILE /var/log/xinetd/qpopper
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
